# dashboard/scrum-api.service
[Unit]
Description=SCRUM API Server - Multi-Agent Coordination Backend
Documentation=https://github.com/tylerbuilds/scrum-mcp
After=network.target

[Service]
Type=simple
User=boski
Group=boski
WorkingDirectory=/mnt/data/projects/scrum-mcp
ExecStart=/home/linuxbrew/.linuxbrew/bin/npx tsx /mnt/data/projects/scrum-mcp/src/index.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scrum-api

# Environment
Environment=NODE_ENV=production
Environment=SCRUM_PORT=4177
Environment=SCRUM_DB_PATH=/mnt/data/projects/scrum-mcp/.scrum/scrum.sqlite
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/boski

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/mnt/data/projects/scrum-mcp

[Install]
WantedBy=multi-user.target
