# dashboard/scrum-frontend.service
[Unit]
Description=SCRUM Frontend - React Dashboard
Documentation=https://github.com/tylerbuilds/scrum-mcp
After=network.target scrum-api.service
Requires=scrum-api.service

[Service]
Type=simple
User=boski
Group=boski
WorkingDirectory=/mnt/data/projects/scrum-mcp/frontend
ExecStartPre=/home/linuxbrew/.linuxbrew/bin/npm run build
ExecStart=/home/linuxbrew/.linuxbrew/bin/npm run preview -- --host 0.0.0.0 --port 5174
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scrum-frontend

# Environment
Environment=NODE_ENV=production
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/boski

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/mnt/data/projects/scrum-mcp/frontend

[Install]
WantedBy=multi-user.target
